# upnp-things.service
# assuming
#	* there is a upnp-things user this should be run as
#	* upnp-things was npm-installed in the upnp-things directory under /mnt/home1
#	* upnp-things requires mDNS host name resolution
# as root,
# define, notice, start and enable the systemd system upnp-things.service
#	cp /mnt/home1/upnp-things/node_modules/upnp-things/upnp-things.service /etc/systemd/system/
#	systemctl daemon-reload
#	systemctl start upnp-things
#	systemctl enable upnp-things
# follow journaled output
#	journalctl -f -u upnp-things

[Unit]
Documentation=https://www.github.com/rtyle/upnp-things
After=network.target avahi-daemon.service
RequiresMountsFor=/mnt/home1

[Service]
Type=simple
User=upnp-things
ExecStart=/usr/bin/node /mnt/home1/upnp-things/node_modules/upnp-things/index.js --legrand-log=warn --upnp-log=warn
Restart=on-failure

[Install]
WantedBy=multi-user.target
